# app/audiostream.timer
[Unit]
Description=Start the audio stream daily

[Timer]
OnCalendar=*-*-* 6:00:00
Persistent=true

[Install]
WantedBy=timers.target